# root/etc/services.d/h2/run
#!/bin/bash
exec s6-setuidgid user /runH2.sh

